# cli/deployment_configs/CAC/dcp-root/etc/xinetd.d/dcp-miner-v8
service dcp-miner-v8
{
	type			= UNLISTED
	port			= 9000
	only_from		= localhost
	flags			= REUSE
	server			= /global/home/tmp10011/dcp-root/opt/dcp/bin/dcp-miner-v8
	server_args		= /global/home/tmp10011/dcp-root/opt/dcp/bin/dcp-miner-control.js -b /global/home/tmp10011/dcp-root/opt/dcp/libexec -d0
        instances               = 80
        socket_type             = stream
        wait                    = no
        protocol                = tcp
        user                    = tmp10011
        group                   = tmpgroup
        cps                     = 5000 1
	per_source		= UNLIMITED
	max_load		= 10
	log_type        	= FILE /global/home/tmp10011/dcp-root/var/log/dcp-miner-v8.log
	log_on_failure		= HOST ATTEMPT
	log_on_success		= HOST PID DURATION EXIT
}

